# scripts/com.boomcatcher.autopush.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<!--
  TEMPLATE ONLY — do NOT load this file directly.
  The old version hardcoded "/Users/mac/..." which breaks on other accounts.
  Run  `bash scripts/install-launchd.sh`  instead; it generates the real plist
  in ~/Library/LaunchAgents/ with your correct home path filled in.
-->
<dict>
    <key>Label</key>
    <string>com.boomcatcher.autopush</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>__PUSH_SCRIPT__</string>
    </array>
    <key>StartInterval</key>
    <integer>3600</integer>
    <key>RunAtLoad</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/tmp/boom-catcher-push.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/boom-catcher-push.err</string>
</dict>
</plist>
